# content-pipeline/automation/launchd/com.xaid.content-evening-retry30.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.xaid.content-evening-retry30</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Users/anastasia/Projects/xAID/xaid-landing/content-pipeline/automation/retry30.sh</string>
  </array>
  <!-- Every 30 min. No-op unless the evening run failed on a TRANSIENT error
       (network/auth 403 — not a usage limit), which arms retry-30-<date>.marker.
       Self-terminates on success or after MAX_RETRIES (see retry30.sh). -->
  <key>StartInterval</key><integer>1800</integer>
  <key>StandardOutPath</key><string>/Users/anastasia/Projects/xAID/xaid-landing/content-pipeline/automation/state/launchd-evening.log</string>
  <key>StandardErrorPath</key><string>/Users/anastasia/Projects/xAID/xaid-landing/content-pipeline/automation/state/launchd-evening.log</string>
</dict>
</plist>
