# content-pipeline/automation/launchd/com.xaid.content-evening.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.xaid.content-evening</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Users/anastasia/Projects/xAID/xaid-landing/content-pipeline/automation/evening.sh</string>
  </array>
  <!-- 18:00 by the Mac's local clock (founder decision 2026-07-07: her local
       time, not CET). Missed runs (Mac asleep) fire on wake. -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>18</integer>
    <key>Minute</key><integer>0</integer>
  </dict>
  <key>StandardOutPath</key><string>/Users/anastasia/Projects/xAID/xaid-landing/content-pipeline/automation/state/launchd-evening.log</string>
  <key>StandardErrorPath</key><string>/Users/anastasia/Projects/xAID/xaid-landing/content-pipeline/automation/state/launchd-evening.log</string>
</dict>
</plist>
